# build-scripts/appliance/template/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid USER_ID
setgid GROUP_ID
chdir /var/www/SITENAME/source

exec ../virtualenv/bin/gunicorn \
    --bind unix:/tmp/SITENAME.socket \
    --access-logfile ../logs/gunicorn_access.log \
    --error-logfile ../logs/gunicorn_error.log \
    django_project.wsgi:application
